# stage2/03-install-node-and-app/files/raspberry-lights.service
[Unit]
Description=Raspberry Light Service
After=network.target

[Service]
WorkingDirectory=/home/pi/raspberry-lights-service
ExecStart=/usr/local/bin/npm start
Restart=on-failure
User=pi
#Environment=PORT=8000

[Install]
WantedBy=multi-user.target
